# curator4-5903045888.conf
[program:curator4-5903045888]
command=python3 /home/curator4/curator4-5903045888/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator4/curator4-5903045888/test.err.log
stdout_logfile=/home/curator4/curator4-5903045888/test.out.log